# auto-zapret-daemon.service
[Unit]
Description=Auto-Zapret Monitor Daemon
After=network.target zapret.service
Wants=zapret.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/auto-zapret
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 -m autozapret.main start
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auto-zapret-daemon

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
